# crontab
*/10 * * * * curl -X POST http://localhost:3000/api/v1/chat/send-chat-history -H "Authorization: Bearer $GARBYX_API_KEY" >> /var/log/cron.log 2>&1

# An empty line is required at the end of this file for a valid cron file.
